# modules/example/.config/profile.d/example.zsh
#!/usr/bin/env zsh

# stub file so there's always something that lives in profile.d
